# custom/files/etc/profile.d/banner.sh
#!/bin/sh
# /etc/profile.d/banner.sh

# Load OpenWrt release info
. /etc/openwrt_release

# ANSI color codes
RED="\033[1;31m"
GREEN="\033[1;32m"
YELLOW="\033[1;33m"
CYAN="\033[1;36m"
WHITE="\033[1;37m"
RESET="\033[0m"

# Big ASCII Art Banner
echo -e "${CYAN}"
cat <<'EOF'
_____/\\\\\\\\\\\____/\\\______________/\\\________/\\\__/\\\\\\\\\\\\________/\\\\\\\\\\\\__/\\\\\\\\\\\\\\\_        
 ___/\\\/////////\\\_\/\\\_____________\/\\\_______\/\\\_\/\\\////////\\\____/\\\//////////__\/\\\///////////__       
  __\//\\\______\///__\/\\\_____________\/\\\_______\/\\\_\/\\\______\//\\\__/\\\_____________\/\\\_____________      
   ___\////\\\_________\/\\\_____________\/\\\_______\/\\\_\/\\\_______\/\\\_\/\\\____/\\\\\\\_\/\\\\\\\\\\\_____     
    ______\////\\\______\/\\\_____________\/\\\_______\/\\\_\/\\\_______\/\\\_\/\\\___\/////\\\_\/\\\///////______    
     _________\////\\\___\/\\\_____________\/\\\_______\/\\\_\/\\\_______\/\\\_\/\\\_______\/\\\_\/\\\_____________   
      __/\\\______\//\\\__\/\\\_____________\//\\\______/\\\__\/\\\_______/\\\__\/\\\_______\/\\\_\/\\\_____________  
       _\///\\\\\\\\\\\/___\/\\\\\\\\\\\\\\\__\///\\\\\\\\\/___\/\\\\\\\\\\\\/___\//\\\\\\\\\\\\/__\/\\\\\\\\\\\\\\\_ 
        ___\///////////_____\///////////////_____\/////////_____\////////////______\////////////____\///////////////__
EOF

# Info below the ASCII
echo -e "${YELLOW}"
echo "-----------------------------------------------------"
echo -e " Version: ${WHITE}${DISTRIB_REVISION}${YELLOW}"
echo -e " Target: ${WHITE}${DISTRIB_TARGET}${YELLOW}"
echo    "-----------------------------------------------------"
echo -e "${RESET}"
